# config/sv/mailsync/run
#!/bin/sh -e
exec 2>&1
exec /usr/bin/snooze -H'*' -M/1 /usr/local/bin/mailsync
