# config/sv/boxsync/run
#!/bin/sh -e
exec 2>&1
exec /usr/bin/snooze -H'*' -M/10 ~/.local/bin/scripts/boxsync.sh 
